# Docker/crontab
* *  *   *   *     /home/Whydah-SPAProxyService/update-service.sh >> /var/log/cron.log 2>&1 && /home/Whydah-SPAProxyService/stop-service.sh >> /var/log/cron.log 2>&1
